# extra/lchatd.service.in
[Unit]
Description=Local Chat dispatcher server

[Service]
ExecStart=@sbindir@/lchatd

[Install]
WantedBy=multi-user.target
